# crontab
0 * * * * /usr/local/bin/backup-mongodb.sh >> /tmp/cron-log 2>&1